# config/supervisor/ODR-encoders.conf
# Radio Monaco
[program:odr-padencoder-defee58a-f768-4ca9-a0f7-88f08b2e8748]
command=/usr/local/bin/odr-padenc
 --dir=/home/odr/config/mot/P01/slide/
 --sleep=0
 --dls=/home/odr/config/mot/P01/dls.dls
 --output=defee58a-f768-4ca9-a0f7-88f08b2e8748
 --label=12
 --label-ins=1200
autostart=true
autorestart=false
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-padencoder-defee58a-f768-4ca9-a0f7-88f08b2e8748.log

# Radio Monaco
[program:odr-audioencoder-defee58a-f768-4ca9-a0f7-88f08b2e8748]
command=/usr/local/bin/odr-audioenc
 --vlc-uri=http://radiomonaco.ice.infomaniak.ch/radiomonaco-128.mp3
 --drift-comp
 --silence=180
 --bitrate=88
 --rate=48000
 --channels=2
 --aaclc
 --pad=34
 --pad-socket=defee58a-f768-4ca9-a0f7-88f08b2e8748
 --write-icy-text=/home/odr/config/mot/P01/dls.dls
 -e tcp://127.0.0.1:9001
 --stats=/var/tmp/P01.stats
autostart=true
autorestart=false
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-audioencoder-defee58a-f768-4ca9-a0f7-88f08b2e8748.log

# Capital FM
[program:odr-padencoder-7158c498-1ef1-4580-b385-23a3c06ba985]
command=/usr/local/bin/odr-padenc
 --dir=/home/odr/config/mot/P02/slide
 --sleep=0
 --dls=/home/odr/config/mot/P02/dls.dls
 --output=7158c498-1ef1-4580-b385-23a3c06ba985
 --label=12
 --label-ins=1200
autostart=true
autorestart=false
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-padencoder-7158c498-1ef1-4580-b385-23a3c06ba985.log

# Capital FM
[program:odr-audioencoder-7158c498-1ef1-4580-b385-23a3c06ba985]
command=/usr/local/bin/odr-audioenc
 --gst-uri=http://media-ice.musicradio.com/CapitalMP3
 --drift-comp
 --silence=60
 --bitrate=88
 --rate=48000
 --channels=2
 --aaclc
 --pad=34
 --pad-socket=7158c498-1ef1-4580-b385-23a3c06ba985
 --write-icy-text=/home/odr/config/mot/P02/dls.dls
 -e tcp://127.0.0.1:9002
 --stats=/var/tmp/P02.stats
autostart=true
autorestart=false
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-audioencoder-7158c498-1ef1-4580-b385-23a3c06ba985.log

# Swiss Pop
[program:odr-padencoder-f93a9393-bc99-48a5-afa9-c681c7a91066]
command=/usr/local/bin/odr-padenc
 --dir=/home/odr/config/mot/P03/slide/
 --sleep=0
 --dls=/home/odr/config/mot/P03/dls.dls
 --output=f93a9393-bc99-48a5-afa9-c681c7a91066
 --label=12
 --label-ins=1200
autostart=true
autorestart=true
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-padencoder-f93a9393-bc99-48a5-afa9-c681c7a91066.log

# Swiss Pop
[program:odr-audioencoder-f93a9393-bc99-48a5-afa9-c681c7a91066]
command=/usr/local/bin/odr-audioenc
 --gst-uri=https://stream.srg-ssr.ch/m/rsp/mp3_128
 --drift-comp
 --silence=60
 --bitrate=88
 --rate=48000
 --channels=2
 --sbr
 --pad=34
 --pad-socket=f93a9393-bc99-48a5-afa9-c681c7a91066
 --write-icy-text=/home/odr/config/mot/P03/dls.dls
 -e tcp://127.0.0.1:9003
 --stats=/var/tmp/f93a9393-bc99-48a5-afa9-c681c7a91066.stats
autostart=true
autorestart=true
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-audioencoder-f93a9393-bc99-48a5-afa9-c681c7a91066.log

# NRJ+
[program:odr-padencoder-c353716e-f276-4e9c-b5e8-782018edd78e]
command=/usr/local/bin/odr-padenc
 --dir=/home/odr/config/mot/P04/slide/
 --sleep=0
 --dls=/home/odr/config/mot/P04/dls.dls
 --output=c353716e-f276-4e9c-b5e8-782018edd78e
 --label=12
 --label-ins=1200
autostart=true
autorestart=true
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-padencoder-c353716e-f276-4e9c-b5e8-782018edd78e.log

# NRJ+
[program:odr-audioencoder-c353716e-f276-4e9c-b5e8-782018edd78e]
command=/usr/local/bin/odr-audioenc
 --gst-uri=https://stream.rcs.revma.com/1a6hdnzbebuvv
 --drift-comp
 --silence=60
 --bitrate=88
 --rate=48000
 --channels=2
 --sbr
 --pad=48
 --pad-socket=c353716e-f276-4e9c-b5e8-782018edd78e
 --write-icy-text=/home/odr/config/mot/P04/dls.dls
 -e tcp://127.0.0.1:9004
 --stats=/var/tmp/c353716e-f276-4e9c-b5e8-782018edd78e.stats
autostart=true
autorestart=true
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-audioencoder-c353716e-f276-4e9c-b5e8-782018edd78e.log

# Radio Paradise
[program:odr-padencoder-9020a318-7bed-448f-9ac6-e0558a74feeb]
command=/usr/local/bin/odr-padenc
 --dir=/home/odr/config/mot/P05/slide/
 --sleep=0
 --dls=/home/odr/config/mot/P05/dls.dls
 --output=9020a318-7bed-448f-9ac6-e0558a74feeb
 --label=12
 --label-ins=1200
autostart=true
autorestart=true
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-padencoder-9020a318-7bed-448f-9ac6-e0558a74feeb.log

# Radio Paradise
[program:odr-audioencoder-9020a318-7bed-448f-9ac6-e0558a74feeb]
command=/usr/local/bin/odr-audioenc
 --gst-uri=https://stream.radioparadise.com/world-etc-32
 --drift-comp
 --silence=60
 --bitrate=88
 --rate=48000
 --channels=2
 --sbr
 --pad=34
 --pad-socket=9020a318-7bed-448f-9ac6-e0558a74feeb
 --write-icy-text=/home/odr/config/mot/P05/dls.dls
 -e tcp://127.0.0.1:9005
 --stats=/var/tmp/9020a318-7bed-448f-9ac6-e0558a74feeb.stats
autostart=true
autorestart=true
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-audioencoder-9020a318-7bed-448f-9ac6-e0558a74feeb.log

# KINK
[program:odr-padencoder-3b06ac10-5ecb-4f1a-9136-7d03907af4da]
command=/usr/local/bin/odr-padenc
 --dir=/home/odr/config/mot/P06/slide
 --sleep=0
 --dls=/home/odr/config/mot/P06/dls.dls
 --output=3b06ac10-5ecb-4f1a-9136-7d03907af4da
 --label=12
 --label-ins=1200
autostart=true
autorestart=true
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-padencoder-3b06ac10-5ecb-4f1a-9136-7d03907af4da.log

# KINK
[program:odr-audioencoder-3b06ac10-5ecb-4f1a-9136-7d03907af4da]
command=/usr/local/bin/odr-audioenc
 --gst-uri=https://22663.live.streamtheworld.com/KINK_SC
 --drift-comp
 --silence=60
 --bitrate=88
 --rate=48000
 --channels=2
 --sbr
 --pad=34
 --pad-socket=3b06ac10-5ecb-4f1a-9136-7d03907af4da
 --write-icy-text=/home/odr/config/mot/P06/dls.dls
 -e tcp://127.0.0.1:9006
 --stats=/var/tmp/3b06ac10-5ecb-4f1a-9136-7d03907af4da.stats
autostart=true
autorestart=true
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-audioencoder-3b06ac10-5ecb-4f1a-9136-7d03907af4da.log

# Maxi 80
[program:odr-padencoder-f5634b8a-9aea-4a3f-8a7e-b7bfca12a004]
command=/usr/local/bin/odr-padenc
 --dir=/home/odr/config/mot/P07/slide/
 --sleep=0
 --dls=/home/odr/config/mot/P07/dls.dls
 --output=f5634b8a-9aea-4a3f-8a7e-b7bfca12a004
 --label=12
 --label-ins=1200
autostart=true
autorestart=true
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-padencoder-f5634b8a-9aea-4a3f-8a7e-b7bfca12a004.log

# Maxi 80
[program:odr-audioencoder-f5634b8a-9aea-4a3f-8a7e-b7bfca12a004]
command=/usr/local/bin/odr-audioenc
 --gst-uri=https://audio1.maxi80.com/;?1707061581287
 --drift-comp
 --silence=180
 --bitrate=88
 --rate=48000
 --channels=2
 --sbr
 --pad=34
 --pad-socket=f5634b8a-9aea-4a3f-8a7e-b7bfca12a004
 --write-icy-text=/home/odr/config/mot/P07/dls.dls
 -e tcp://127.0.0.1:9007
 --stats=/var/tmp/f5634b8a-9aea-4a3f-8a7e-b7bfca12a004.stats
autostart=true
autorestart=true
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-audioencoder-f5634b8a-9aea-4a3f-8a7e-b7bfca12a004.log

# Radio Meuh
[program:odr-padencoder-2f258770-f550-4117-8f94-6acb157bc1db]
command=/usr/local/bin/odr-padenc
 --dir=/home/odr/config/mot/P08/slide/
 --sleep=0
 --dls=/home/odr/config/mot/P08/dls.dls
 --output=2f258770-f550-4117-8f94-6acb157bc1db
 --label=12
 --label-ins=1200
autostart=true
autorestart=true
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-padencoder-2f258770-f550-4117-8f94-6acb157bc1db.log

# Radio Meuh
[program:odr-audioencoder-2f258770-f550-4117-8f94-6acb157bc1db]
command=/usr/local/bin/odr-audioenc
 --gst-uri=https://radiomeuh2.ice.infomaniak.ch/radiomeuh2-128.mp3
 --drift-comp
 --silence=180
 --bitrate=88
 --rate=48000
 --channels=2
 --sbr
 --pad=34
 --pad-socket=2f258770-f550-4117-8f94-6acb157bc1db
 --write-icy-text=/home/odr/config/mot/P08/dls.dls
 -e tcp://127.0.0.1:9008
 --stats=/var/tmp/2f258770-f550-4117-8f94-6acb157bc1db.stats
autostart=true
autorestart=true
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-audioencoder-2f258770-f550-4117-8f94-6acb157bc1db.log

# FG Antwerpen
[program:odr-padencoder-758a8c51-97f8-48f7-953e-dfa7cbe9ebcf]
command=/usr/local/bin/odr-padenc
 --dir=/home/odr/config/mot/P09/slide/
 --sleep=0
 --dls=/home/odr/config/mot/P09/dls.dls
 --output=758a8c51-97f8-48f7-953e-dfa7cbe9ebcf
 --label=12
 --label-ins=1200
autostart=true
autorestart=true
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-padencoder-758a8c51-97f8-48f7-953e-dfa7cbe9ebcf.log

# FG Antwerpen
[program:odr-audioencoder-758a8c51-97f8-48f7-953e-dfa7cbe9ebcf]
command=/usr/local/bin/odr-audioenc
 --gst-uri=https://n07a-eu.rcs.revma.com/5wesqhfap98uv?rj-ttl=5&rj-tok=AAABjXTuLCQAcnBsxlO9uU79Pw
 --drift-comp
 --silence=180
 --bitrate=88
 --rate=48000
 --channels=2
 --sbr
 --pad=34
 --pad-socket=758a8c51-97f8-48f7-953e-dfa7cbe9ebcf
 --write-icy-text=/home/odr/config/mot/P09/dls.dls
 -e tcp://127.0.0.1:9009
 --stats=/var/tmp/758a8c51-97f8-48f7-953e-dfa7cbe9ebcf.stats
autostart=true
autorestart=true
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-audioencoder-758a8c51-97f8-48f7-953e-dfa7cbe9ebcf.log
